# Sample8_ML.Flask.Web.Application/aws_eb_deploy/cron/aws_update.cron

# Update AWS IP ranges daily at midnight
0 0 * * * root . /etc/environment; /usr/local/bin/aws_update_ips.sh >> /proc/1/fd/1 2>> /proc/1/fd/2



