# Sample8_ML.Flask.Web.Application/aws_eb_deploy/cron/aws_update.cron

# Update AWS IP ranges daily at midnight
0 0 * * * /usr/local/bin/aws_update_ips.sh >> /dev/stdout 2>> /dev/stderr


